# devcontainer/devcontainer-pod.container
# Quadlet container file for systemd integration
# Place this file in ~/.config/containers/systemd/
# Then run: systemctl --user daemon-reload && systemctl --user start devcontainer-pod

[Unit]
Description=Development Container for VS Code
After=local-fs.target

[Container]
# Container image (build first with instructions in README)
Image=localhost/devcontainer:latest

# Container name
ContainerName=devcontainer

# Mount workspace from host
Volume=%h/devmachine:/home/developer

# Mount SSH keys (read-only for security)
Volume=%h/.ssh:/home/developer/.ssh:ro

# Mount Windows SSH keys (for Positron/Windows IDE access)
Volume=/mnt/c/Users/rsimon/.ssh:/home/developer/.ssh-windows:ro

# Timezone (inherit from host)
Volume=/etc/localtime:/etc/localtime:ro

# Environment variables
Environment=TERM=xterm-256color

# Expose SSH port (container 22 -> host 2222)
PublishPort=2222:22

# Map host UID/GID to same UID/GID inside the container
# Without this, rootless Podman remaps host user → container root,
# making mounted volumes unwritable by the container's non-root user
UserNS=keep-id

# WSL compatibility - disable problematic features
Notify=false
PodmanArgs=--cgroups=no-conmon

[Service]
# Restart policy
Restart=always
RestartSec=10

# Give time for graceful shutdown
TimeoutStopSec=30

[Install]
WantedBy=default.target
